# ops/heuesta-event-reminder.timer
[Unit]
Description=Run HEU ESTA event reminders hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
