# ops/heuesta-backup.timer
[Unit]
Description=Run HEU ESTA backup daily at 03:30

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
